# config/rc.local
#!/bin/sh -e
# rc.local for BOSS
sudo python /home/pi/Worker.py
exit 0